# systemd/fb-chrome-lean.service
[Unit]
Description=fb Chrome backend using your real profile (Default) — cookies + OpenCLI
After=fb-opencli.service
Wants=fb-opencli.service
ConditionPathExists=/usr/bin/google-chrome
# ConditionPathExists is filled loosely; install may adjust
ConditionPathExists=%h/.config/google-chrome/Default

[Service]
Type=simple
ExecStart=%h/bin/fb-chrome-lean
Environment=HOME=%h
Environment=FB_CHROME_USER_DATA=%h/.config/google-chrome
Environment=FB_CHROME_PROFILE_DIR=Default
# attach = never spawn Chrome; wait for your normal Chrome (no new windows)
Environment=FB_CHROME_MODE=attach
Restart=on-failure
RestartSec=5
MemoryMax=2G
MemoryHigh=1536M
Nice=5
KillMode=control-group
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fb-chrome-lean

[Install]
WantedBy=default.target
